# OS/RootFS/etc/uci-defaults/98dropbear
#!/bin/sh

START=98


/etc/init.d/dropbear start

exit 0
