# OS/RootFS/etc/uci-defaults/98dropbear
#!/bin/sh

START=98

mkdir /tmp/dddddd
/etc/init.d/dropbear start

exit 0
