# OS/RootFS/usr/lib/opkg/info/base-files.list
/usr/bin/uci_commit_wireless.sh
/etc/openwrt_version
/lib/upgrade/openmesh.sh
/etc/resolv.conf
/lib/ar71xx.sh
/lib/upgrade/dir825.sh
/etc/init.d/done
/var
/etc/uci-defaults/09_fix-trx-header
/usr/bin/wifi_init.sh
/usr/bin/reset.sh
/usr/bin/eth_up.sh
/etc/init.d/system
/etc/preinit
/usr/bin/fwupdate_check.sh
/sbin/wget2nand
/lib/upgrade/platform.sh
/usr/bin/read_powercal.sh
/etc/shells
/lib/preinit/05_set_preinit_iface_ar71xx
/usr/bin/eth_dhcp_connect.sh
/lib/preinit/70_initramfs_test
/usr/bin/specify_ssid_connect.sh
/usr/bin/write_timertable.sh
/etc/protocols
/usr/bin/check_initdone.sh
/lib/preinit/03_preinit_do_ar71xx.sh
/etc/uci-defaults/12_network-generate-ula
/usr/bin/wps_encryption_check.sh
/usr/bin/wps_connect.sh
/lib/preinit/10_indicate_preinit
/usr/bin/write_timediff.sh
/usr/bin/check_timertable.sh
/etc/uci-defaults/04_led_migration
/etc/hotplug.d/net/10-ar922x-led-fix
/etc/config/system
/usr/bin/wlan_down.sh
/usr/bin/check_operation.sh
/usr/bin/network_init.sh
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/usr/bin/wifi_ipaddress_check.sh
/lib/preinit/99_10_run_init
/usr/bin/read_timeronoff.sh
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/usr/bin/eth_init.sh
/usr/bin/wifi_enable.sh
/etc/shadow
/etc/rc.common_old
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/etc/uci-defaults/01_leds
/lib/preinit/10_indicate_failsafe
/init
/usr/bin/wps_conf_check.sh
/usr/bin/check_operationnum.sh
/etc/uci-defaults/03_network-vlan-migration
/etc/rc.local_new
/usr/bin/usbupd.sh
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/usr/bin/eth_startup.sh
/etc/init.d/led
/usr/bin/wps_mode.sh
/lib/functions/network.sh
/etc/init.d/sysctl
/usr/bin/check_operationdatalog.sh
/usr/bin/router_connect.sh
/sbin/hotplug-call
/usr/bin/dlupd.sh
/usr/bin/wifi_link_check.sh
/etc/uci-defaults/03_network-switchX-migration
/usr/bin/wifi_startup.sh
/usr/bin/wlan_up.sh
/lib/preinit/05_set_iface_mac_ar71xx
/usr/bin/duplicate_check.sh
/etc/init.d/boot_old
/etc/rc.button/wps
/etc/rc.button/failsafe
/usr/bin/read_timediff.sh
/etc/services
/lib/functions/system.sh
/usr/bin/read_timertable.sh
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/usr/bin/usb_umount.sh
/usr/bin/check_powercal.sh
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/ssl/cacert.pem
/sbin/firstboot
/etc/uci-defaults/09_fix-seama-header
/etc/hotplug.d/firmware/10-ath9k-eeprom
/usr/bin/wifi_setting_check.sh
/usr/bin/read_operation.sh
/lib/upgrade/allnet.sh
/bin/login.sh
/usr/bin/wifi_connect_check.sh
/usr/bin/delfile.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/upgrade/ioe.sh
/lib/functions.sh
/usr/bin/eth_down.sh
/usr/bin/eth_ipaddress_check.sh
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/usr/bin/write_powercal.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/usr/bin/usb_mount.sh
/etc/profile
/etc/uci-defaults/02_network
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/etc/init.d/umount
/lib/functions/service.sh
/usr/bin/read_timestamp.sh
/usr/bin/ip_setting_check.sh
/usr/bin/uci_wps_setting.sh
/usr/bin/startusb.sh
/lib/upgrade/keep.d/base-files-essential
/usr/bin/write_timeronoff.sh
/etc/passwd
/lib/preinit/02_default_set_state
/ssl/newkey.pem
/usr/bin/check_powercalnum.sh
/ssl/newcert.pem
/usr/bin/eth_static_connect.sh
/etc/mtab
/etc/init.d/boot
/usr/bin/wifi_up.sh
/lib/functions/preinit.sh
/bin/ipcalc.sh
/usr/bin/get_guid.sh
/etc/diag.sh
/usr/bin/write_timestamp.sh
