# OS/RootFS/usr/lib/opkg/info/firewall.list
/etc/init.d/firewall
/sbin/fw3
/etc/config/firewall
/etc/firewall.user
/etc/hotplug.d/iface/20-firewall
