# OS/RootFS/usr/lib/opkg/info/fwupdate.list
/etc/rc.d/S09fw
/usr/sbin/fwupdate
/usr/sbin/fwtool
