# OS/RootFS/usr/lib/opkg/info/mcproxy.list
/etc/hotplug.d/iface/50-mcproxy
/etc/config/mcproxy
/etc/init.d/mcproxy
/usr/sbin/mcproxy
