# OS/RootFS/usr/lib/opkg/info/uhttpd.list
/usr/sbin/uhttpd
/etc/init.d/uhttpd
/etc/config/uhttpd
/lib/upgrade/keep.d/uhttpd
